# luci-app-webdav/root/etc/init.d/webdav
#!/bin/sh /etc/rc.common
# Copyright (C) 2025 Lean

START=99
STOP=50
USE_PROCD=1

PROG="/usr/bin/webdav"
CONF="/var/etc/webdav.yaml"

gen_conf() {
	mkdir -p /var/etc

	cat <<-EOF > $CONF
# Auto generated by /etc/init.d/webdav
address: 0.0.0.0
port: ${port}
auth: true
tls: false
cert: cert.pem
key: key.pem
prefix: /

cors:
  enabled: true
  credentials: true
  allowed_headers:
    - Depth
  allowed_hosts:
    - http://localhost:${port}
  allowed_methods:
    - GET
    - PUT
    - POST
    - DELETE
    - PROPFIND
  exposed_headers:
    - Content-Length
    - Content-Range

scope: .
modify: true
rules: []

users:
  - username: ${username}
    password: ${password}
    scope: ${scope}
    modify: true
EOF
}

start_service() {
	config_load webdav

	config_get port    config port    8000
	config_get username config username
	config_get password config password
	config_get scope   config scope   "/"
	config_get_bool enable config enable 0

	[ "$enable" -eq 1 ] || return

	gen_conf

	procd_open_instance
	procd_set_param command $PROG -c $CONF
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_set_param respawn
	procd_close_instance
}

service_triggers() {
	procd_add_reload_trigger "webdav"
}
